# supervisor.conf
[supervisord]
nodaemon=true

[program:npm rent-a-car]
command=npm run start:prod

[program:npm compodocs]
command=npm run doc:serve